# package/blueprint/scripts/fleet.gem.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>fleet.gem</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/opt/homebrew/bin/gemini</string>
    <string>--yolo</string>
    <string>-p</string>
    <string>Run your heartbeat protocol. Read ~/fleet/MISSION_CONTROL.md first.</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/miguelrodriguez/fleet/gem</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Minute</key><integer>0</integer></dict>
    <dict><key>Minute</key><integer>10</integer></dict>
    <dict><key>Minute</key><integer>20</integer></dict>
    <dict><key>Minute</key><integer>30</integer></dict>
    <dict><key>Minute</key><integer>40</integer></dict>
    <dict><key>Minute</key><integer>50</integer></dict>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/miguelrodriguez/fleet/logs/gem.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/miguelrodriguez/fleet/logs/gem_err.log</string>
</dict>
</plist>
